# systemd/leds.service
[Unit]
Description=LED lightings

[Service]
ExecStart=/usr/bin/python "/home/jheka/raspberrypi/gpio01.py"

[Install]
WantedBy=default.target

